# layers/meta-yummy/recipes-app/yumy-app/files/yumy-app.service
[Unit]
Description=The god of Yumi

[Service]
ExecStart=/opt/yumy-app/MyApp --fb

[Install]
WantedBy=multi-user.target